# other/shadowtls/shadowtls.service.template
[Unit]
Description=ss-server faketls
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/hiddify-config/other/shadowtls
ExecStart=/opt/hiddify-config/other/shadowtls/shadowtls
Environment=RUST_LOG=error
Environment=CONFIG_FILE=/opt/hiddify-config/other/shadowtls/shadowtls.conf
Restart=on-failure
StandardOutput=file:/opt/hiddify-config/log/system/ssplugin-shadowtls.out.log
StandardError=file:/opt/hiddify-config/log/system/ssplugin-shadowtls.err.log
[Install]
WantedBy=multi-user.target
