# shadowsocks/ss-v2ray.service.template
[Unit]
Description=ss-server v2ray
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=ss-server -c /opt/GITHUB_REPOSITORY/shadowsocks/config-v2ray.json
Restart=on-failure

[Install]
WantedBy=multi-user.target
